# dnspod_ddns.service
[Unit]
Description=update hexingzhi.com DNS.

[Service]
ExecStart=/root/publish/dnspod_ddns -t your_token -d google.com
